# packaging/simeis.service
[Unit]
Description=Simeis Media Server
After=network.target

[Service]
Type=simple
User=simeis
ExecStart=/usr/bin/simeis
Restart=on-failure

[Install]
WantedBy=multi-user.target